# scripts/config/systemd/user/profile-offload.service
[Unit]
Description=Move unneeded profile files

[Service]
Type=oneshot
ExecStart=/bin/sh -c '~/.local/share/profdiv.sh'

[Install]
WantedBy=default.target
